# daemon/habpack-map.conf
#!upstart

description "Upstart script for habpack-map"

start on startup
stop on shutdown

respawn
respawn limit 10 5

script
    exec node /srv/susf/daemon/main.js
end script
